# root/etc/services.d/syslogd/run
#!/usr/bin/with-contenv bash

/sbin/syslogd