# root/etc/periodic/hourly/get_iplayer_pvr
#!/usr/bin/with-contenv bash

exec \
    s6-svc -o /var/run/s6/services/get_iplayer
